# config/com.claude.autorenew.sleepaware.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.autorenew.sleepaware</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/yashmehta/src/cc-auto-ping/scripts/claude-auto-renew-sleepaware.sh</string>
    </array>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>TZ</key>
        <string>Asia/Tokyo</string>
    </dict>
    
    <key>StartCalendarInterval</key>
    <array>
        <!-- Monday 7:00 AM JST -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- Monday 12:01 PM JST -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>1</integer>
        </dict>
        <!-- Tuesday 7:00 AM JST -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- Tuesday 12:01 PM JST -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>1</integer>
        </dict>
        <!-- Wednesday 7:00 AM JST -->
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- Wednesday 12:01 PM JST -->
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>1</integer>
        </dict>
        <!-- Thursday 7:00 AM JST -->
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- Thursday 12:01 PM JST -->
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>1</integer>
        </dict>
        <!-- Friday 7:00 AM JST -->
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- Friday 12:01 PM JST -->
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>1</integer>
        </dict>
        <!-- Saturday 7:00 AM JST -->
        <dict>
            <key>Weekday</key>
            <integer>6</integer>
            <key>Hour</key>
            <integer>7</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- Saturday 12:01 PM JST -->
        <dict>
            <key>Weekday</key>
            <integer>6</integer>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>1</integer>
        </dict>
    </array>
    
    <key>StandardOutPath</key>
    <string>/Users/yashmehta/src/cc-auto-ping/logs/claude-autorenew-sleepaware.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/yashmehta/src/cc-auto-ping/logs/claude-autorenew-sleepaware-error.log</string>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>